# supervisord.conf
[supervisord]
nodaemon=true
user=www-data

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
user=www-data

[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
user=www-data
